# roles/media/files/mnt-media.automount
[Unit]
Description = nfs mount for media files

[Automount]
What=srv01.griend.dev:/mnt/pool0/Media
Where=/mnt/media
Type=nfs
Options=defaults,ro
TimeoutSec=30

[Install]
WantedBy=multi-user.target

